# crontab
*/30 8-23 * * * python /main.py >> /logfile

